# daemon/qubes-db.service.d/30_qubes-kvm.conf
[Unit]
ConditionPathExists=/var/run/qubes-service/hypervisor-kvm
After=kvmchand.service
Requires=kvmchand.service
